# main.cron
@reboot /bin/bash /u01/aipoc/No.1-SQL-Assist/main.sh
